# server/backup/passengers-db-backup.service
[Unit]
Description=Passengers backend DB backup
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/passengers-backend/ops/passengers-db-backup.sh
